# plugins/color/csd-color.target.in
[Unit]
Description=GNOME Color management

BindsTo=csd-color.service
After=csd-color.service

Requisite=gnome-session-initialized.target
After=gnome-session-initialized.target
PartOf=gnome-session-initialized.target
Before=gnome-session.target
